# roles/os_hardening/templates/etc/profile.d/tmout.sh.j2
# Logout Timeout
export TMOUT={{ os_security_auto_logout }}
readonly TMOUT
